# c#/projects/MultithreadingExam/MultithreadingExam/copied-files/5ffea6126f02e78b9099eb4614d2d339f03ca5a8.xml - copy.xml
<?xml version="1.0" encoding="utf-16"?>
<Task xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>$(@%systemRoot%\system32\usosvc.dll,-103)</Author>
    <Source>$(@%systemRoot%\system32\usosvc.dll,-104)</Source>
    <Description>$(@%systemRoot%\system32\usosvc.dll,-105)</Description>
    <URI>\Microsoft\Windows\UpdateOrchestrator\Schedule Scan Static Task</URI>
    <SecurityDescriptor>D:P(A;;FA;;;SY)(A;;FRFX;;;LS)(A;;FRFX;;;BA)</SecurityDescriptor>
    <Version>1.0</Version>
  </RegistrationInfo>
  <Triggers>
    <!-- trigger WNF_SHEL_OOBE_USER_LOGON_COMPLETE = {0xa3bc2475, 0x0d83063e} -->
    <WnfStateChangeTrigger>
      <StateName>7524bca33e06830d</StateName>
      <Enabled>true</Enabled>
      <Delay>PT2H5M</Delay>
      <Data>01</Data>
      <DataOffset>0</DataOffset>
    </WnfStateChangeTrigger>
    <!-- trigger WNF_DEP_OOBE_COMPLETE = {0xa3bc0c75, 0x41960b29} -->
    <WnfStateChangeTrigger>
      <StateName>750cbca3290b9641</StateName>
      <Enabled>true</Enabled>
      <Delay>PT2H5M</Delay>
      <Data>01</Data>
      <DataOffset>0</DataOffset>
    </WnfStateChangeTrigger>
    <!-- trigger WNF_OLIC_OS_EDITION_CHANGE = {0xa3bc5075, 0x028f0222} -->
    <WnfStateChangeTrigger>
      <StateName>7550bca322028f02</StateName>
      <Enabled>true</Enabled>
    </WnfStateChangeTrigger>
    <!-- trigger WNF_CI_SMODE_CHANGE = {0xa3bc0875, 0x41c6072e} -->
    <WnfStateChangeTrigger>
      <StateName>7508bca32e07c641</StateName>
      <Enabled>true</Enabled>
    </WnfStateChangeTrigger>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription>&lt;QueryList&gt;&lt;Query Id="0" Path="System"&gt;&lt;Select Path="System"&gt;*[System[EventID=8202]]&lt;/Select&gt;&lt;/Query&gt;&lt;/QueryList&gt;</Subscription>
    </EventTrigger>
  </Triggers>
  <Settings>
    <AllowHardTerminate>true</AllowHardTerminate>
    <AllowStartOnDem*******>true</AllowStartOnDem*******>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <DisallowStartOnRemoteAppSession>false</DisallowStartOnRemoteAppSession>
    <Enabled>true</Enabled>
    <ExecutionTimeLimit>PT72H</ExecutionTimeLimit>
    <Hidden>false</Hidden>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <Priority>7</Priority>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <StartWhenAvailable>true</StartWhenAvailable>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <WakeToRun>false</WakeToRun>
  </Settings>
  <Principals>
    <Principal id="LocalSystem">
      <UserId>S-1-5-18</UserId>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Actions Context="LocalSystem">
    <Exec>
      <Comm*******>%systemroot%\system32\usoclient.exe</Comm*******>
      <Arguments>StartScan</Arguments>
    </Exec>
  </Actions>
</Task>
